# debian/ironic-discoverd.init.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          ironic-discoverd
# Required-Start:    $network $local_fs $remote_fs $syslog
# Required-Stop:     $remote_fs
# Should-Start:      postgresql mysql keystone rabbitmq-server ntp
# Should-Stop:       postgresql mysql keystone rabbitmq-server ntp
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: ironic-discoverd
# Description:       ironic-discoverd
### END INIT INFO

# Author: Thomas Goirand <zigo@debian.org>

DESC="OpenStack Ironic Discoverd"
PROJECT_NAME=ironic-discoverd
NAME=${PROJECT_NAME}
